# pocsag_decoder/rootfs/etc/services.d/decode/run
#!/usr/bin/with-contenv bash

BAUD=$(jq -r ".baud" /data/options.json)

exec multimon-ng -t raw -a ${BAUD} -
